# gen-config/supervisor-devpi.conf
[program:devpi-server]
command=/home/merzouka/esi/security/presentation/demo/venv/bin/devpi-server 
priority = 999
startsecs = 5
redirect_stderr = True
autostart = True
